# API/installer/timers/corosync.timer
[Unit]
Description=corosync_timer

[Timer]
OnBootSec=40
Unit=corosync.service

[Install]
WantedBy=multi-user.target
